# .config/systemd/user/xscreensaver.service
[Unit]
Description=XscreenSaver
[Service]
ExecStart=xscreensaver
[Install]
WantedBy=default.target
